# gophernicus@.service
[Unit]
Description=Gophernicus gopher server

[Service]
ExecStart=/usr/sbin/in.gophernicus -h %H
StandardInput=socket
User=nobody
